# resources/gatekeeper.service
[Unit]
Description=Gatekeeper Multi-Port Forwarder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/gatekeeper.sh
Restart=always
RestartSec=2
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
